# lib/systemd/tio-trim-cache.service
[Unit]
Description=Trim TIO's output cache

[Service]
Type=oneshot
User=apache
Group=apache
ExecStart=/srv/bin/trim-cache
