# lib/systemd/tio-trim-cache.timer
[Unit]
Description=Trim TIO's output cache

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
