# lib/systemd/tio-letsrenew.timer
[Unit]
Description=Renews Let's Encrypt SSL certificates, once a week

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
