# lib/systemd/tio-letsrenew.service
[Unit]
Description=Renews Let's Encrypt SSL certificates

[Service]
Type=oneshot
ExecStart=/usr/bin/env certbot renew
