# lib/systemd/tioapi.socket
[Unit]
Description=TIO API socket

[Socket]
ListenStream=/run/tioapi{}.sock
Accept=yes
MaxConnections=1

[Install]
WantedBy=sockets.target
